# src/main/java/study/schedule/jcrontab/crontab
* * * * * com.aweb.test.NumTest
* * * * * com.aweb.test.LetterTest